# packages/module-build-service/module-build-service-1.0.2.tar.gz/module-build-service-1.0.2/conf/mbs-scheduler.service
[Unit]
Description=Scheduler for the Module Build Service
After=network.target
Documentation=https://pagure.io/fm-orchestrator

[Service]
ExecStart=/usr/bin/mbs-daemon
Type=simple
User=fedmsg
Group=fedmsg
Restart=on-failure

[Install]
WantedBy=multi-user.target

